# manifests/crane/files/docker.upstart
# docker

description "docker"

start on runlevel [234]
stop on starting rc RUNLEVEL=[0156]

exec /usr/bin/docker -d
respawn
