# manifests/crane/files/docker_daemon.upstart
# docker

description "docker"

start on runlevel [234]
stop on starting rc RUNLEVEL=[0156]

exec /usr/bin/docker -d -H tcp://127.0.0.1:4243 -H unix:///var/run/docker.sock
respawn
